# scripts/com.stltexturizer.http.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stltexturizer.http</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>@REPO_ROOT@</string>
  <key>ProgramArguments</key>
  <array>
    <string>@PYTHON3@</string>
    <string>-m</string>
    <string>http.server</string>
    <string>8000</string>
    <string>--bind</string>
    <string>0.0.0.0</string>
  </array>
  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/stltexturizer-http.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/stltexturizer-http.err.log</string>
</dict>
</plist>
